# Praktikum 1 Sisop/task-3/crontabs
#automation megang archive sm download
*/2 * * * * ~/cloud_storage/automation.sh

0 */2 * * * ~/cloud_storage/archive.sh
